# etc/profile.d/alias.sh
# /etc/profile.d/alias.sh

alias dd="dd status=progress"
alias df="df -h"
alias ls="ls -hF --group-directories-first"
alias mkdir="mkdir -p"
alias mv="mv -i"
alias rm="rm -i"
