# etc/profile.d/alias.sh
# /etc/profile.d/alias.sh

alias dd="dd status=progress"
alias df="df -h"
alias ls="ls -hF"
alias mkdir="mkdir -p"
alias mv="mv -i"
alias rm="rm -i"
alias rsync="rsync -haczP --ignore-existing"
alias emacsq="emacs -q -l ~/.emacs.d/tini.el"
