# etc/profile.d/defaults.sh
#!/usr/bin/env sh

export LANG="en_US.UTF-8"
export PAGER="less"
export PROMPT_COMMAND="history -a"
