# etc/profile.d/gambit.sh
#!/usr/bin/env sh

if [ -e "/usr/local/gambit" ]; then
    export PATH="$PATH:/usr/local/gambit/current/bin"
fi
