# etc/profile.d/gentoo.sh
#!/usr/bin/env sh

if [ -e /etc/profile.env ]; then
    . /etc/profile.env
fi
