# etc/profile.d/gentoo.sh
# /etc/profile.d/gentoo.sh

if [ -e /etc/profile.env ]; then
    . /etc/profile.env
fi
