# etc/profile.d/lang.sh
#!/usr/bin/env sh

CHARSET="${CHARSET:-\"\"}"
LANG="${LANG:-en_US.UTF-8}"
LC_TYPE="${LC_TYPE:-\"\"}"

export CHARSET LANG LC_TYPE
