# etc/profile.d/python.sh
# /etc/profile.d/python.sh

alias pip="pip --disable-pip-version-check"
