# etc/profile.d/rust.sh
#!/usr/bin/env sh

if [ -f "$HOME/.cargo/env" ]; then
    . ~/.cargo/env
fi
