# etc/profile.d/st.sh
#!/usr/bin/env sh

if [ "$TERM" == "st-256color" ]; then
    export TERM="xterm-256color"
fi
