# etc/profile.d/vim.sh
# /etc/profile.d/vim.sh

if [ -x $(command -v vimx) ]; then
    alias vim="vimx"
fi
